# templates/default/solr.conf.erb
# /etc/init/solr.conf

author "David Tee"
description "Starts Solr server"

respawn
start on startup solr_push

script
	echo "Starting solr server..."
	cd <%= @path %> && java -Djava.util.logging.config.file=<%= @config_dir %>/logging.properties -Dsolr.solr.home=<%= @config_dir %> -jar start.jar >> <%= @log_dir %>/solr-console.log 2>&1
end script
